# Antminer-L3-v1.02-20180114_0-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/avahi-daemon.postinst
if [ -z "$D" ]; then
	killall -q -HUP dbus-daemon || true
fi
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	systemctl $OPTS enable avahi-daemon.service

	if [ -z "$D" -a "enable" = "enable" ]; then
		systemctl restart avahi-daemon.service
	fi
fi
