# Antminer-L3-v1.02-20180114_0-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/initscripts.list
/etc/user_defined_lcd.factory
/etc/default/devpts
/etc/init.d/sysfs.sh
/etc/init.d/umountfs
/etc/rc0.d/S30urandom
/sbin/minerlcd
/etc/rcS.d/S39ntpdate.sh
/etc/init.d/read-only-rootfs-hook.sh
/etc/init.d/checkroot.sh
/etc/rcS.d/S37bitmainer_setup.sh
/etc/rc0.d/S20sendsigs
/etc/rc0.d/S40umountfs
/etc/cgminer.conf.factory
/etc/rcS.d/S38devpts.sh
/etc/rcS.d/S45mountnfs.sh
/etc/rc2.d/S99rmnologin.sh
/etc/default/volatiles/00_core
/etc/init.d/banner.sh
/etc/user_setting.factory
/etc/rcS.d/S29read-only-rootfs-hook.sh
/etc/init.d/populate-volatile.sh
/etc/rc4.d/S99rmnologin.sh
/etc/rcS.d/S30urandom
/etc/init.d/bootmisc.sh
/etc/rcS.d/S06alignment.sh
/etc/rcS.d/S02banner.sh
/etc/rc6.d/S31umountnfs.sh
/etc/init.d/devpts.sh
/etc/init.d/urandom
/etc/rc0.d/S90halt
/etc/init.d/ntpdate.sh
/usr/bin/compile_time
/sbin/monitorcg
/etc/shadow.factory
/etc/rc0.d/S25save-rtc.sh
/etc/init.d/rmnologin.sh
/etc/rcS.d/S40ntpd
/var/lib/urandom/
/etc/rc6.d/S30urandom
/etc/init.d/halt
/etc/init.d/hostname.sh
/etc/init.d/sendsigs
/etc/rcS.d/S37populate-volatile.sh
/etc/init.d/save-rtc.sh
/etc/network.conf.factory
/etc/init.d/network.sh
/etc/rc6.d/S25save-rtc.sh
/sbin/pic.txt
/etc/init.d/cgminer.sh
/etc/rc3.d/S99rmnologin.sh
/etc/rcS.d/S10checkroot.sh
/etc/rc6.d/S20sendsigs
/etc/rc6.d/S40umountfs
/etc/rcS.d/S39hostname.sh
/etc/rc6.d/S90reboot
/etc/init.d/pgnand.sh
/etc/init.d/reboot
/etc/rcS.d/S35mountall.sh
/etc/init.d/mountdevtmpfs.sh
/etc/init.d/alignment.sh
/etc/rc5.d/S99rmnologin.sh
/etc/device_table
/etc/rcS.d/S02sysfs.sh
/etc/init.d/bitmainer_setup.sh
/etc/init.d/umountnfs.sh
/etc/rcS.d/S55bootmisc.sh
/etc/rcS.d/S70cgminer.sh
/etc/init.d/mountnfs.sh
/etc/init.d/functions.initscripts
/etc/rcS.d/S02mountdevtmpfs.sh
/etc/rcS.d/S38network.sh
/sbin/reset.sh
/etc/init.d/mountall.sh
/etc/rc0.d/S31umountnfs.sh
/etc/init.d/single
