# Antminer-L3-v1.02-20180114_0-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/ntp.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	systemctl $OPTS enable ntpd.service

	if [ -z "$D" -a "enable" = "enable" ]; then
		systemctl restart ntpd.service
	fi
fi
