# Antminer-L3-v1.02-20180114_0-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/ntp.prerm
#!/bin/sh
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop ntpd.service
	fi

	systemctl disable ntpd.service
fi
