# Antminer-L3-v1.02-20180114_0-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/ntpdate.list
/usr/bin/ntpdate-sync
/etc/network/if-up.d/ntpdate-sync
/etc/network/if-up.d/
/etc/network/
/etc/default/ntpdate
/lib/systemd/system/ntpdate.service
/usr/bin/ntpdate
