# Antminer-D3-v2.06-20180105_1_beta-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/avahi-daemon.prerm
#!/bin/sh
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop avahi-daemon.service
	fi

	systemctl disable avahi-daemon.service
fi
