# Antminer-D3-v2.06-20180105_1_beta-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/busybox-syslog.list
/lib/systemd/
/lib/systemd/system/
/etc/init.d/
/etc/syslog-startup.conf.busybox
/lib/systemd/system/busybox-syslog.service
/etc/default/busybox-syslog
/etc/default/
/lib/systemd/system/syslog.service
/lib/systemd/system/busybox-klogd.service
/etc/init.d/syslog.busybox
