# Antminer-D3-v2.06-20180105_1_beta-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/lighttpd.postinst
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
if type update-rc.d >/dev/null 2>/dev/null; then
	update-rc.d $OPT lighttpd defaults 70
fi
