# Antminer-D3-v2.06-20180105_1_beta-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/lighttpd.list
/www/pages/resources/icons/Thumbs.db
/www/pages/type/titilliumweb-extralightitalic-webfont.woff
/www/pages/cgi-bin/upgrade_clear.cgi
/www/pages/resources/icons/vlan.png
/www/pages/resources/cbi/reload.gif
/www/pages/resources/icons/wifi_big.png
/www/pages/resources/cbi/file.gif
/www/pages/js/
/www/pages/resources/icons/signal-25-50.png
/www/pages/network.html
/www/pages/css/ie7.css
/www/pages/cgi-bin/get_system_info.cgi
/www/pages/minerStatus.html
/www/pages/cgi-bin/network_diag.cgi
/www/pages/type/titilliumweb-semibold-webfont.woff
/www/pages/css/cascade.css
/www/pages/resources/icons/encryption.png
/www/pages/resources/icons/loading.gif
/www/pages/resources/
/www/pages/resources/cbi/help.gif
/www/pages/resources/cbi/fieldadd.gif
/www/pages/resources/icons/ethernet_disabled.png
/www/pages/resources/icons/signal-0.png
/www/pages/resources/icons/port_down.png
/www/pages/css/ie8.css
/www/pages/type/titilliumweb-extralight-webfont.woff
/www/pages/resources/icons/wifi_big_disabled.png
/www/pages/networkDiagnostics.html
/www/pages/cgi-bin/get_miner_status.cgi
/www/pages/type/titilliumweb-italic-webfont.eot
/www/pages/cgi-bin/get_kernel_log.cgi
/www/pages/type/titilliumweb-extralight-webfont.eot
/etc/lighttpd-htdigest.user
/www/pages/resources/icons/wifi_disabled.png
/www/pages/type/titilliumweb-black-webfont.eot
/www/pages/resources/icons/port_up.png
/www/pages/resources/icons/tunnel_disabled.png
/www/pages/type/titilliumweb-semibolditalic-webfont.eot
/www/pages/css/
/www/pages/type/titilliumweb-black-webfont.woff
/www/pages/type/titilliumweb-extralightitalic-demo.html
/www/pages/type/titilliumweb-bolditalic-demo.html
/www/pages/administration.html
/www/pages/resources/cbi/find.gif
/www/pages/type/titilliumweb-regular-webfont.woff
/www/pages/cgi-bin/kill_cgminer.cgi
/www/pages/type/titilliumweb-lightitalic-webfont.eot
/www/pages/css/style.css
/www/pages/resources/cbi/
/www/pages/dav/
/www/pages/monitor.html
/etc/init.d/lighttpd
/www/pages/cgi-bin/reset_conf.cgi
/www/pages/cgi-bin/set_miner_conf.cgi
/www/pages/cgi-bin/get_network_info.cgi
/www/pages/resources/icons/ethernet.png
/www/pages/resources/icons/bridge_disabled.png
/www/pages/type/titilliumweb-italic-webfont.ttf
/www/pages/type/titilliumweb-italic-webfont.svg
/www/pages/resources/cbi/edit.gif
/www/pages/cgi-bin/cgi_lib.cgi
/www/pages/type/titilliumweb-extralight-webfont.ttf
/www/pages/resources/cbi/add.gif
/www/pages/type/titilliumweb-black-webfont.ttf
/www/pages/type/titilliumweb-extralight-webfont.svg
/www/pages/type/titilliumweb-black-webfont.svg
/etc/rcS.d/S60lighttpd
/www/pages/resources/icons/switch.png
/www/pages/type/titilliumweb-semibolditalic-webfont.ttf
/www/pages/type/titilliumweb-semibolditalic-webfont.svg
/www/logs/
/www/pages/resources/icons/encryption_disabled.png
/www/pages/css/tabbg.png
/www/pages/upgrade.html
/www/pages/reboot.html
/www/pages/cgi-bin/reboot.cgi
/www/pages/resources/icons/signal-0-25.png
/www/pages/type/titilliumweb-extralightitalic-webfont.eot
/www/pages/type/titilliumweb-lightitalic-webfont.ttf
/www/pages/kernelLog.html
/www/pages/type/titilliumweb-light-demo.html
/www/pages/type/titilliumweb-lightitalic-webfont.svg
/www/pages/cgi-bin/create_conf_backup.cgi
/www/pages/type/titilliumweb-lightitalic-demo.html
/www/pages/type/titilliumweb-black-demo.html
/www/pages/cgi-bin/set_network_conf.cgi
/www/pages/images/antminer_logo.png
/www/pages/resources/cbi/link.gif
/www/pages/cgi-bin/minerAdvanced.cgi
/www/pages/type/titilliumweb-regular-demo.html
/www/pages/type/titilliumweb-extralightitalic-webfont.ttf
/www/pages/type/titilliumweb-extralight-demo.html
/www/pages/type/titilliumweb-extralightitalic-webfont.svg
/etc/lighttpd.conf
/www/pages/type/titilliumweb-semibolditalic-webfont.woff
/www/pages/resources/cbi/folder.gif
/usr/sbin/lighttpd
/www/pages/type/
/www/pages/type/titilliumweb-bolditalic-webfont.woff
/www/pages/resources/icons/bridge.png
/www/pages/cgi-bin/passwd.cgi
/www/pages/type/titilliumweb-semibold-webfont.eot
/www/pages/images/
/www/pages/type/titilliumweb-bold-demo.html
/www/pages/resources/icons/switch_disabled.png
/www/pages/resources/icons/
/www/pages/cgi-bin/minerConfiguration.cgi
/www/pages/cgi-bin/
/www/pages/resources/cbi/reset.gif
/www/pages/type/titilliumweb-semibold-demo.html
/www/pages/type/titilliumweb-regular-webfont.eot
/www/pages/type/type.css
/www/pages/type/titilliumweb-bold-webfont.eot
/www/pages/js/json2.min.js
/www/pages/resources/cbi/down.gif
/www/pages/js/jquery-1.10.2.js
/usr/sbin/lighttpd-angel
/www/pages/type/titilliumweb-semibold-webfont.ttf
/www/pages/type/titilliumweb-semibold-webfont.svg
/www/pages/cgi-bin/monitor.cgi
/www/pages/cgi-bin/kill_bmminer.cgi
/www/pages/resources/icons/signal-50-75.png
/www/pages/css/grid.css
/www/pages/cgi-bin/miner_lcd.sh
/www/
/www/pages/type/titilliumweb-bolditalic-webfont.eot
/www/pages/type/titilliumweb-regular-webfont.ttf
/www/pages/type/titilliumweb-regular-webfont.svg
/www/pages/cgi-bin/upgrade.cgi
/www/pages/type/titilliumweb-bold-webfont.ttf
/www/pages/resources/cbi/remove.gif
/www/pages/type/titilliumweb-bold-webfont.svg
/www/pages/cgi-bin/minerStatus.cgi
/www/pages/cgi-bin/minerStatus.html
/www/pages/cgi-bin/upload_conf.cgi
/www/pages/type/titilliumweb-light-webfont.woff
/www/var/
/www/pages/type/titilliumweb-bold-webfont.woff
/www/pages/cgi-bin/get_status_api.cgi
/www/pages/resources/icons/signal-75-100.png
/www/pages/resources/cbi/user.gif
/www/pages/type/titilliumweb-bolditalic-webfont.ttf
/www/pages/resources/cbi/up.gif
/www/pages/type/titilliumweb-bolditalic-webfont.svg
/www/pages/minerConfiguration.html
/www/pages/type/titilliumweb-lightitalic-webfont.woff
/www/pages/type/titilliumweb-italic-demo.html
/www/pages/js/xhr.js
/www/pages/resources/icons/wifi.png
/www/pages/index.html
/www/pages/resources/cbi/download.gif
/www/pages/type/titilliumweb-light-webfont.eot
/www/pages/resources/cbi/apply.gif
/www/pages/resources/icons/tunnel.png
/www/pages/resources/icons/signal-none.png
/www/pages/cgi-bin/get_miner_conf.cgi
/www/pages/type/titilliumweb-semibolditalic-demo.html
/www/pages/resources/cbi/arrow.gif
/www/pages/
/www/pages/css/ie6.css
/www/pages/type/titilliumweb-light-webfont.ttf
/www/pages/resources/icons/vlan_disabled.png
/www/pages/type/titilliumweb-light-webfont.svg
/www/pages/resources/cbi/save.gif
/www/pages/type/titilliumweb-italic-webfont.woff
/www/pages/resources/cbi/key.gif
