# Antminer-D3-v2.06-20180105_1_beta-blissz/_initramfs.bin.SD.extracted/var/lib/opkg/info/ntp.list
/usr/bin/ntpd
/etc/ntp.conf
/etc/init.d/ntpd
/lib/systemd/system/ntpd.service
