# Antminer-D3-201711201715-0M/_initramfs.bin.SD.extracted/var/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/usr/sbin/dropbearmulti
/etc/pam.d/
/etc/dropbear/
/etc/init.d/dropbear
/usr/sbin/dropbearkey
/usr/sbin/dropbearconvert
/usr/sbin/dropbear
/etc/pam.d/dropbear
/etc/rcS.d/S40dropbear
