# Antminer-D3-201711201715-0M/_initramfs.bin.SD.extracted/var/lib/opkg/info/monitor-recobtn.list
/etc/rcS.d/S70monitor-recobtn
/usr/bin/monitor-recobtn
/usr/bin/factory_config_reset.sh
/etc/init.d/monitor-recobtn
