# Antminer-D3-201711201715-0M/_initramfs.bin.SD.extracted/var/lib/opkg/info/monitor-ipsig.list
/usr/bin/monitor-ipsig
/etc/rcS.d/S100monitor-ipsig.sh
/etc/init.d/monitor-ipsig.sh
