# Antminer-D3-201711201715-0M/_initramfs.bin.SD.extracted/var/lib/opkg/info/ntpdate.prerm
#!/bin/sh
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop ntpdate.service
	fi

	systemctl disable ntpdate.service
fi
